# rootfs/etc/s6/php/run
#!/bin/bash

if [ -f ./setup ]
then
  source ./setup
fi

if [ -f ./configure-custom ]
then
  source ./configure-custom
fi

echo "> starting php"
exec /usr/sbin/php-fpm7
